# myflaskproject-deb/DEBIAN/postinst
#!/bin/bash
# servisi başlat
systemctl enable myapp.service
systemctl start myapp.service
