# myflaskproject-deb/DEBIAN/prerm
#!/bin/bash
# servisi durdur
systemctl stop myapp.service
systemctl disable myapp.service
